# user-services/.config/systemd/user/power-monitor.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=/home/bezer/scripts/battery_script.sh

[Install]
WantedBy=default.target
